# roles/clamav/templates/clamav.freshclam.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>clamav.freshclam</string>
	<key>ProgramArguments</key>
	<array>
		<string>{{ brew_prefix }}/bin/freshclam</string>
		<string>--daemon</string>
		<string>--foreground</string>
		<string>--config-file</string>
		<string>{{ clamav_freshclam_config_file }}</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>{{ clamav_freshclam_stderr_log_file }}</string>
</dict>
</plist>